# ops/staging/systemd/levelchannel-staging.service
# LevelChannel — staging instance.
#
# Twin of the production `levelchannel.service` (lives in the private
# runbook) but with:
#   - distinct port (3001 vs prod 3000)
#   - distinct env file (/etc/levelchannel-staging/env)
#   - distinct work-dir (/srv/levelchannel-staging vs prod /srv/levelchannel)
#   - distinct user (`levelchannel-staging` vs prod `levelchannel`)
#   - distinct journald tag
#
# Same hardening posture as prod (SystemCallFilter, ProtectSystem, etc.)
# so seccomp regression tests cover both. See
# `scripts/seccomp-regression-check.sh`.
#
# Operator install:
#   cp ops/staging/systemd/levelchannel-staging.service \
#      /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now levelchannel-staging.service
#
# Full activation runbook: docs/staging-setup.md

[Unit]
Description=LevelChannel Next.js — staging
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=levelchannel-staging
Group=levelchannel-staging
WorkingDirectory=/srv/levelchannel-staging
EnvironmentFile=/etc/levelchannel-staging/env
# Runtime env file written by `ops/staging/scripts/autodeploy-staging.sh`
# before each `systemctl start`. Carries dynamic per-deploy values
# (currently just GIT_SHA, surfaced by /api/health.version and used by
# .github/workflows/staging-e2e.yml to poll for the promoted SHA).
# Leading `-` makes systemd tolerate the file being absent on the
# very first boot before autodeploy has run.
EnvironmentFile=-/var/lib/levelchannel-staging/runtime.env
Environment="NODE_ENV=production"
Environment="LC_ENV=staging"
Environment="PORT=3001"
# Path is relative to WorkingDirectory (/srv/levelchannel-staging).
# Resolves to /srv/levelchannel-staging/node_modules/.bin/next at
# runtime — installed by `npm ci` during autodeploy.
ExecStart=/srv/levelchannel-staging/node_modules/.bin/next start --port 3001
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=levelchannel-staging

# --- Hardening ---
# MemoryDenyWriteExecute is intentionally NOT set: Node 20 V8 in
# baseline JIT mode mprotects allocated pages W↔X at runtime (verified
# crash on 2026-06-03 during the first staging cold-start). The
# pkey_* allowlist + @system-service filter alone are tolerant of the
# JIT path and still catch the @privileged / @resources regression
# class that scripts/seccomp-regression-check.sh asserts on.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_NETLINK is required for libuv's uv_interface_addresses (Linux
# getifaddrs uses rtnetlink under the hood). Without it Next 16's
# start-server.js post-listen callback throws ERR_SYSTEM_ERROR 97
# (EAFNOSUPPORT) and never registers its request handler. Verified
# 2026-06-03 against /api/health curl from 127.0.0.1.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# Explicit pkey_* allowlist for Node + V8 + glibc parity.
# scripts/seccomp-regression-check.sh enforces this list per-unit.
SystemCallFilter=pkey_alloc pkey_mprotect pkey_free
ReadWritePaths=/srv/levelchannel-staging /var/log/levelchannel-staging
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
