# scripts/systemd/levelchannel-stale-orders.service
# Reference systemd unit for scripts/cancel-stale-orders.mjs.
# `scripts/activate-prod-ops.sh` renders the placeholders below during
# installation from the private runbook.

[Unit]
Description=LevelChannel hourly stale-pending-order cleanup
After=network.target

[Service]
Type=oneshot
User=levelchannel
Group=levelchannel
WorkingDirectory=__LEVELCHANNEL_APP_DIR__
EnvironmentFile=__LEVELCHANNEL_ENV_FILE__
ExecStart=/usr/bin/node scripts/cancel-stale-orders.mjs
StandardOutput=journal
StandardError=journal

# Codex 2026-05-08 (Wave 8 #6 LOW) — least-privilege sandboxing.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
# MemoryDenyWriteExecute=true — incompatible with Node.js V8 JIT (kernel
# kills the process via SIGSYS / status 31 because V8 needs W+X pages
# for compiled code). Verified live on prod 2026-05-08: smoke run with
# this directive aborted with `code=dumped, status=31/SYS`. Removed; the
# remaining 13 directives still give meaningful defense-in-depth.
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @debug @mount @cpu-emulation @obsolete
