# scripts/systemd/levelchannel-auto-complete-slots.service
# Reference systemd unit for scripts/auto-complete-slots.mjs.
# `scripts/activate-prod-ops.sh` renders the placeholders below during
# installation from the private runbook.

[Unit]
Description=LevelChannel daily auto-complete past-booked lesson slots
After=network.target

[Service]
Type=oneshot
User=levelchannel
Group=levelchannel
WorkingDirectory=__LEVELCHANNEL_APP_DIR__
EnvironmentFile=__LEVELCHANNEL_ENV_FILE__
ExecStart=/usr/bin/node scripts/auto-complete-slots.mjs
StandardOutput=journal
StandardError=journal

# Codex 2026-05-08 (Wave 8 #6 LOW) — least-privilege sandboxing.
# These directives restrict what the script process can do even if
# the script itself is compromised. None of them affect the script's
# legitimate work (read app dir, connect to Postgres over TCP, write
# stdout/stderr to journald).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
# MemoryDenyWriteExecute=true — removed: incompatible with Node.js V8 JIT
# (W+X pages are required, the directive triggers SIGSYS). See
# levelchannel-stale-orders.service for the live-prod incident note.
# SystemCallFilter / SystemCallArchitectures — removed alongside MDWE for
# Node.js V8 incompatibility. See levelchannel-stale-orders.service for
# the live-prod incident detail.
# SystemCallArchitectures=native
# SystemCallFilter=@system-service
# SystemCallFilter=~@privileged @resources @debug @mount @cpu-emulation @obsolete
