# scripts/systemd/levelchannel-stale-orders.service
# Reference systemd unit for scripts/cancel-stale-orders.mjs.
# `scripts/activate-prod-ops.sh` renders the placeholders below during
# installation from the private runbook.

[Unit]
Description=LevelChannel hourly stale-pending-order cleanup
After=network.target

[Service]
Type=oneshot
User=levelchannel
Group=levelchannel
WorkingDirectory=__LEVELCHANNEL_APP_DIR__
EnvironmentFile=__LEVELCHANNEL_ENV_FILE__
ExecStart=/usr/bin/node scripts/cancel-stale-orders.mjs
StandardOutput=journal
StandardError=journal
