# scripts/systemd/levelchannel-webhook-flow-alert.service
# Reference systemd unit for the webhook-flow alert probe. Copy to
# /etc/systemd/system/ on the VPS during activation; see OPERATIONS.md
# §9 "Webhook-flow alerting" for the full procedure.
#
# This unit is one-shot — invoked by the .timer companion every 30 min.
# All env (DATABASE_URL, RESEND_API_KEY, EMAIL_FROM, ALERT_EMAIL_TO)
# is read from /etc/levelchannel.env so we don't duplicate secrets.

[Unit]
Description=LevelChannel webhook-flow alert probe
After=network.target

[Service]
Type=oneshot
User=levelchannel
Group=levelchannel
WorkingDirectory=/var/www/levelchannel
EnvironmentFile=/etc/levelchannel.env
ExecStart=/usr/bin/node scripts/webhook-flow-alert.mjs
StandardOutput=journal
StandardError=journal
